# deploy/fitlab.service
[Unit]
Description=Fitlab Gunicorn
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/fitlab
EnvironmentFile=/var/www/fitlab/.env
ExecStart=/var/www/fitlab/.venv/bin/gunicorn fitlab_project.wsgi:application \
    --bind 127.0.0.1:__FITLAB_GUNICORN_PORT__ \
    --workers 3 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile -
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
